# files/default/automacchanger.service
[Unit]
Description=Auto MAC Changer

[Service]
User=root
ExecStart=/usr/local/bin/automacchanger

[Install]
WantedBy=multi-user.target
